# services/importer/run
#!/bin/bash
export RAILS_ENV=RAILSENV
export PATH=DAEMON_PATH
export GEM_HOME=GEMHOME
export GEM_PATH=GEMPATH
exec 2>&1
exec LOGJAM_DIR/vendor/plugins/logjam/script/import_amqp_log IMPORTER
