# services/importer/run
#!/bin/sh
export PATH=RUBY_PATH:/usr/local/bin:$PATH
exec LOGJAM_DIR/vendor/plugins/logjam/script/import_amqp_log IMPORTER
