# services/parser/run
#!/bin/sh
export RAILS_ENV=RAILSENV
export PATH=RUBY_PATH:/usr/local/bin:$PATH
exec LOGJAM_DIR/vendor/plugins/logjam/script/parse_amqp_log IMPORTER
