# services/proxy/run
#!/bin/bash
exec 2>&1
exec logjam-proxy LOGJAM_SERVICE_TARGET_DIR/proxy.conf
