# services/cimporter/run
#!/bin/bash
# ulimit -n 16384
exec 2>&1
exec /usr/local/bin/logjam-importer -c LOGJAM_SERVICE_TARGET_DIR/cimporter.conf "$@"
